# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql --datadir=/var/lib/mysql
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/mariadb-stdout.log
stderr_logfile=/var/log/mariadb-stderr.log

[program:php-fpm]
command=/bin/sh -c "exec /usr/sbin/php-fpm%(ENV_PHP_VERSION)s --nodaemonize --fpm-config /etc/php/%(ENV_PHP_VERSION)s/fpm/php-fpm.conf"
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/php-fpm-stdout.log
stderr_logfile=/var/log/php-fpm-stderr.log

[program:caddy]
command=caddy run --config %(ENV_DEVCONTAINER_DIR)s/Caddyfile
environment=DOMAIN="%(ENV_DOMAIN)s",WORKSPACE_ROOT="%(ENV_WORKSPACE_ROOT)s",WEB_ROOT="%(ENV_WEB_ROOT)s",PHP_FPM_SOCKET="%(ENV_PHP_FPM_SOCKET)s",HTTP_PORT="%(ENV_HTTP_PORT)s",HTTPS_PORT="%(ENV_HTTPS_PORT)s",DEVCONTAINER_DIR="%(ENV_DEVCONTAINER_DIR)s"
autostart=true
autorestart=true
priority=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
